# modules/install-as-service/srkbz-infra.service
[Unit]
Description=srkbz infra
After=network.target

[Service]
Type=oneshot
ExecStart=${SRKBZ_INFRA_ROOT}/taskw
WorkingDirectory=${SRKBZ_INFRA_ROOT}
