# modules/infra-service/srkbz-infra.path
[Path]
PathExists=/srv/srkbz/triggers/APPLY

[Install]
WantedBy=multi-user.target
