# modules/ebrow-service/ebrow.path
[Path]
PathChanged=${TRIGGER_HOME}/RUN

[Install]
WantedBy=multi-user.target
