# modules/taskw-service/taskw.path
[Path]
PathExists=${TASKW_ROOT}/.state/taskw/triggers/RUN

[Install]
WantedBy=multi-user.target
