# modules/taskw-service/taskw.path
[Path]
PathChanged=/opt/taskw/triggers/RUN

[Install]
WantedBy=multi-user.target
